# luci-app-adguardhome/root/etc/uci-defaults/40_luci-adguardhome
#!/bin/sh

# Migrate the historical default without overriding user-selected paths.
configpath="$(uci -q get AdGuardHome.AdGuardHome.configpath)"
if [ "$configpath" = /etc/config/AdGuardHome.yaml ] && [ ! -e "$configpath" ]; then
	uci -q set AdGuardHome.AdGuardHome.configpath=/etc/AdGuardHome.yaml
	configpath=/etc/AdGuardHome.yaml
fi
[ -n "$configpath" ] || configpath=/etc/AdGuardHome.yaml
if [ ! -e "$configpath" ]; then
	parent="${configpath%/*}"; [ -n "$parent" ] || parent=/
	mkdir -p "$parent"
	cp /usr/share/AdGuardHome/AdGuardHome_template.yaml "$configpath"
	chmod 600 "$configpath"
fi
uci -q delete AdGuardHome.AdGuardHome.ucitracktest
uci -q commit AdGuardHome

if [ "$(uci -q get firewall.AdGuardHome)" = include ] && \
	[ "$(uci -q get firewall.AdGuardHome.path)" = /usr/share/AdGuardHome/firewall.start ]; then
	uci -q delete firewall.AdGuardHome
	uci -q commit firewall
fi

rm -f /tmp/luci-indexcache
exit 0
