# luci-app-ddns-go/root/etc/uci-defaults/99-luci-ddns-go
#!/bin/sh

chmod 755 /usr/share/rpcd/ucode/luci.ddns-go
rm -f /tmp/luci-indexcache /tmp/luci-modulecache/*
[ -n "$IPKG_INSTROOT" ] || /etc/init.d/rpcd restart
exit 0
